# scripts/daily-server-restart.service
[Unit]
Description=Restart Server

[Service]
ExecStart=/bin/sh -c "cd /root/mycoda && docker compose restart"
